# services.conf
[program:sshd]
command=/usr/sbin/sshd -D

[program:php5-fpm]
command=/usr/sbin/php5-fpm -F -c /etc/php5/fpm

[program:mysqld]
command=/usr/sbin/mysqld --skip-innodb-use-native-aio --basedir=/usr --datadir=/var/lib/mysql --plugin-dir=/usr/lib/mysql/plugin --user=mysql --log-error=/var/log/mysql/error.log --pid-file=/var/run/mysqld/mysqld.pid --socket=/var/run/mysqld/mysqld.sock --port=3306

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"

